# crontab
PATH=/usr/digital_library_server

19 0-23 * * *  ./backup.py -h
19 19 1-31 * * ./backup.py -d
19 19 * * =4   ./backup.py -w