# debian/cron.d
PATH=~/digital_library

19 0-23 * * *  root ./backuper.sh h
19 19 1-31 * * root ./backuper.sh d
19 19 * * =4   root ./backuperroot .sh w